# daemon/systemd/faduinod.service
[Unit]
Description=faduinod

[Service]
Type=simple
WorkingDirectory=/home/kajuha/catkin_ws/src/faduino/daemon/systemd
ExecStart=/home/kajuha/catkin_ws/src/faduino/daemon/systemd/faduinod.sh
User=root
Group=root
Restart=always

[Install]
WantedBy=multi-user.target
